# services/nbawinspool.service.d/input.conf
[Service]
WorkingDirectory=path/to/deploy/dir
Environment="GITHUB_WEBHOOK_SECRET_TOKEN=<Github secret token here>"
